# scripts/systemd/fpdf.service
[Unit]
Description=fpdf PDF filling server
After=network.target
StartLimitIntervalSec=100
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/a_star/home/hcma/work/ai/fpdf
# Use the wrapper script which handles node/pnpm setup
ExecStart=/a_star/home/hcma/work/ai/fpdf/scripts/fpdf --listen-all --port 8002
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
