# etc/systemd/fpdf.service
[Unit]
Description=fpdf PDF filling server
After=network.target
StartLimitIntervalSec=100
StartLimitBurst=5

[Service]
Type=simple
# Assumes fpdf is located in ~/work/ai/fpdf
# Adjust the path below if your repository is in a different location
WorkingDirectory=@@REPO_DIR@@
ExecStart=@@REPO_DIR@@/scripts/fpdf --listen-all --port 8002
Restart=always
RestartSec=5
StandardOutput=append:%h/scratch/fpdf/fpdf.log
StandardError=append:%h/scratch/fpdf/fpdf.log
ExecStartPost=/usr/bin/curl --silent --fail --max-time 5 --retry 12 --retry-delay 1 --retry-connrefused http://127.0.0.1:8002/health

[Install]
WantedBy=default.target
